# config/pro/web.service
# /etc/systemd/system/web.service

[Service]
ExecStart=/home/cotown/www/webs/gen.sh