# config/pro/web.path
# /etc/systemd/system/web.path

[Path]
PathChanged=/home/cotown/app/web/

[Install]
WantedBy=multi-user.target